# rtorrent/.config/systemd/user/ipv4-blocklist-updater.timer
[Unit]
Description=Update rtorrent's IPv4 Blocklist daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
